# deploy/systemd/mr-sentinel.timer
[Unit]
Description=Run mr-sentinel poller every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=false

[Install]
WantedBy=timers.target
